# target/linux/rockchip/armv8/base-files/etc/init.d/rockchip-fan
#!/bin/sh /etc/rc.common

START=21
start() {
    /usr/bin/rockchip-pwm-fan.sh &
}
